# dist/install.sh
#!/bin/bash
set -eux

# Fail fast if we're not on OS X >= 10.6.0. (Snippets taken from Pow.)
if [ "$(uname -s)" != "Darwin" ]; then
  echo "Sorry, this installation script requires Mac OS X to run." >&2
  exit 1
elif [ "$(expr "$(sw_vers -productVersion | cut -f 2 -d .)" \>= 6)" = 0 ]; then
  echo "This script requires Mac OS X 10.6 or later." >&2
  exit 1
fi

if [ -z "$GOPATH" ]; then
  echo "Please install Go first. You can try 'brew install golang && export GOPATH=\"$HOME/go\"'" >&2
  exit 1
fi

mkdir -p "$HOME/.pow"
# save PATH so we can reuse it later
echo "$PATH" > "$HOME/.pow/.path"

go install github.com/jonasschneider/gow

sudo tee /Library/LaunchDaemons/com.jonasschneider.gow.firewall.plist > /dev/null <<END
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jonasschneider.gow.firewall</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-c</string>
    <string>echo "rdr pass inet proto tcp from any to self port {80,20559} -> 127.0.0.1 port 20559" | pfctl -a "com.apple/250.PowFirewall" -f - -E</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>UserName</key>
  <string>root</string>
</dict>
</plist>
END

sudo launchctl load -Fw /Library/LaunchDaemons/com.jonasschneider.gow.firewall.plist 2>/dev/null

sudo tee /etc/resolver/dev > /dev/null <<END
# Generated by Gow. <3 Pow!
nameserver 127.0.0.1
port 20560
END

agent_plist="$HOME/Library/LaunchAgents/com.jonasschneider.gow.gowd.plist"

cat > $agent_plist <<END
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>Label</key>
        <string>com.jonasschneider.gow.gowd</string>
        <key>ProgramArguments</key>
        <array>
                <string>$HOME/.pow/.run</string>
        </array>
        <key>KeepAlive</key>
        <true/>
        <key>RunAtLoad</key>
        <true/>
</dict>
</plist>
END

cat > "$HOME/.pow/.run" <<END
#!/bin/sh
exec $GOPATH/bin/gow > "$HOME/Library/Logs/gowd.log" 2>&1
END

chmod +x "$HOME/.pow/.run"

echo "*** Starting the Pow server..."
launchctl unload "$agent_plist" 2>/dev/null || true
launchctl load -Fw "$agent_plist" 2>/dev/null

# Use networksetup(8) to create a temporary network location,
# switch to it, switch back to the original location, then
# delete the temporary location. This forces reloading of the
# system network configuration.
function reload_network_configuration() {
  echo "*** Reloading system network configuration..."
  local location=$(networksetup -getcurrentlocation)
  networksetup -createlocation "pow$$" >/dev/null 2>&1
  networksetup -switchtolocation "pow$$" >/dev/null 2>&1
  networksetup -switchtolocation "$location" >/dev/null 2>&1
  networksetup -deletelocation "pow$$" >/dev/null 2>&1
}

reload_network_configuration

echo "*** Installed"
